# archiso/airootfs/etc/skel/.config/autostart/a.sh.desktop
[Desktop Entry]
Exec=/home/santiago/a.sh
Icon=dialog-scripts
Name=a.sh
Path=
Type=Application
X-KDE-AutostartScript=true
